# Embedded/Board/JetsonTX2/AWA6280.service
[Unit]
Description=Xserver Autostart

[Service]
ExecStart=/usr/bin/startx
User=root

[Install]
WantedBy=graphical.target
